# systemd/mnt-backupd.automount
[Unit]
Description=Automontaje del disco de copias de seguridad

[Automount]
Where=/mnt/backupd

[Install]
WantedBy=multi-user.target
